# gps/balloon-gps.service
[Unit]
Description=space balloon GPS
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/home/maker/SpaceBalloon/gps/start.sh
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
