# webcam/balloon-webcam3.service
[Unit]
Description=space balloon webcam 3
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/home/maker/SpaceBalloon/webcam/start.sh 3
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
